# bsp/mach-nitrogen6x/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttymxc1::askfirst:-/bin/sh
